# util/systemd/system/7days-watcher.service
[Unit]
Description=7days Watcher Service
After=network.target

[Service]
Type=simple
User=7days
ExecStart=/usr/bin/expect /home/7days/7days-utils/savethehorde.sh
Restart=on-failure


[Install]
WantedBy=multi-user.target
